# systemd/minecraft-backup.service
[Unit]
Description=Minecraft Server Backup
Requisite=minecraft.service

[Service]
Type=oneshot
ExecStart=/usr/bin/free
